# SE2019/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod = 0777
[supervisord]
logfile=/var/run/supervisord.log ;
logfile_maxbytes=50MB        ;
logfile_backups=10           ;
loglevel=info                ;
pidfile=/var/run/supervisord.pid ;
nodaemon=false               ;
minfds=1024                  ;
minprocs=200                 ;
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:Software__service]
command=uwsgi --ini /root/SE2019/script/uwsgi.ini
directory=/root/SE2019
stdout_logfile=/root/SE2019/log/uwsgi_out.log
stderr_logfile=/root/SE2019/log/uwsgi_err.log
stopsignal=QUIT
autostart=true
autorestart=true
user=root
startsecs=10
stopasgroup=true
killasgroup=true
[supervisord]

[program:Software_websocket]
command=daphne -u /root/SE2019/websocket.sock Django.asgi:application
directory=/root/SE2019
stdout_logfile=/root/SE2019/log/websocket_out.log
stderr_logfile=/root/SE2019/log/websocket_err.log
autostart=true
autorestart=true
user=root
startsecs=10
[supervisord]
